# xsession/services/xsettingsd/run
#!/bin/sh
exec xsettingsd -c $HOME/.config/xsettingsd
